# infra/systemd/mt-rotator-restore-verify.timer
[Unit]
Description=Verify an MT Rotator backup every month

[Timer]
OnCalendar=*-*-01 04:30:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=30m
Unit=mt-rotator-restore-verify.service

[Install]
WantedBy=timers.target
